# cloud/src/rasp-cloud/support-files/init.d/rasp.agent
#!/bin/bash

#
# Init file for rasp-cloud
#
### BEGIN INIT INFO
# Provides:          openrasp
### END INIT INFO

export PATH='/usr/sbin:/usr/bin:/sbin:/bin'
export BINPATH='../..'
export BINFILE='./rasp-cloud'
cd $BINPATH


help() {
    $BINFILE -h
    exit 1
}

version() {
    $BINFILE -version
    exit 1    
}

start() {
    $BINFILE -type=agent -d
    return 0
}

reset() {
    $BINFILE -type=reset
    return 0
}

restart() {
    $BINFILE -s restart
    return 0
}

stop() {
   $BINFILE -s stop
}


case "${1}" in
version)
    version
    ;;
restart)
    restart
    ;;
start)
    start
    ;;
reset)
    reset
    ;;
stop)
    stop
    ;;    
*)
    help
    ;;
esac

